# root/etc/s6-overlay/s6-rc.d/init-mod-universal-trust-kachitoritai-ca-install/run
#!/usr/bin/with-contenv bash

update-ca-certificates